# charles_hub/rootfs/etc/services.d/charles_hub/finish
#!/usr/bin/env bashio

if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "CHARLES Hub exited with code $1, halting add-on"
  echo "$1" > /run/s6-linux-init-container-results/exitcode
  exec /run/s6/basedir/bin/halt
fi

bashio::log.info "CHARLES Hub service finished; restart allowed."
