# charles_hub/rootfs/etc/services.d/charles_hub/run
#!/usr/bin/with-contenv bashio
# Minimal service: logs the persona prompt and exits successfully.

message="$(bashio::config 'persona_prompt' 'You are CHARLES.')"
bashio::log.info "CHARLES persona prompt: ${message}"

# Touch an output file in /share to prove execution
echo "${message}" > /share/charles_persona.txt

exit 0
